# ansible/roles/glucose-display-rpi/files/nightscout_rpi.service.j2
[Unit]
Description=RPI PyQT Nightscout Display

[Service]
ExecStart=/usr/bin/python3 /root/app/main.py
WorkingDirectory=/root/app

[Install]
WantedBy=multi-user.target